# build/rcjameica
#!/bin/sh
### BEGIN INIT INFO
# Provides:     jameica
# Required-Start:   $local_fs $remote_fs $network
# Required-Stop:    $local_fs $remote_fs $network
# Default-Start:    3 5
# Default-Stop:     0 1 2 6
# Short-Description:    Jameica Application Server
# Description:      Starts the Jameica Application Server
### END INIT INFO

#
# Determine the base and follow a runlevel link name.
#
base=${0##*/}
link=${base#*[SK][0-9][0-9]}

JAVA=/usr/lib/java/bin/java
BASEDIR=/opt/jameica
WORKDIR=~/.jameica
PASSWORD=CHANGE_THIS_BEFORE_YOU_START_JAMEICA
PIDFILE=$WORKDIR/jameica.pid

### Please do not edit anything below this line! ########################
#                                                                       #

. /etc/rc.status

rc_reset

return=$rc_done
case "$1" in
    start)
      mkdir -p $WORKDIR
      cd $BASEDIR
      echo -n "Starting Jameica "
      $JAVA -jar jameica-linux.jar -d -f $WORKDIR -p $PASSWORD > /dev/null 2>&1 & PID=$!
      rc_status -v
      echo $PID > $PIDFILE
    ;;
    stop)
     if test -e $PIDFILE ; then
       echo -n "Shutting down Jameica "
       kill `cat $PIDFILE`
       rm $PIDFILE
       rc_status -v
     else
      echo "No $SERVERNAME running"
     fi
    ;;
    reload|restart)
      $0 stop
      sleep 1
      $0 start
    ;;

    status)
      echo -n "Checking for Jameica "
      if test -e $PIDFILE ; then
        echo " seems running"
      else
        echo " not running"
      fi
    ;;
    *)
      echo "Usage: $0 {start|stop|restart|status}"
      exit 1
  ;;
esac
rc_exit

